# scripts/kbv2.service
[Unit]
Description=KBv2 Crypto Knowledgebase
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=agentzero
WorkingDirectory=/home/muham/development/kbv2
Environment=PYTHONPATH=/home/muham/development/kbv2/src
EnvironmentFile=/home/muham/development/kbv2/.env.production

# Use production entry point
ExecStart=/usr/bin/uv run python -m knowledge_base.production
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s

# Resource limits
LimitNOFILE=65535
MemoryMax=4G

[Install]
WantedBy=multi-user.target
